# deploy/imagepromptive-deploy-hook.service
[Unit]
Description=ImagePromptive Gitee Deploy Hook
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/www/image-2-prompt-site
EnvironmentFile=/www/image-2-prompt-site/deploy/hook.env
ExecStart=/usr/bin/env python3 /www/image-2-prompt-site/deploy/hook_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
